# packaging/debian/autostart/f9-talk.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=F9 Talk
Comment=Hold-to-talk dictation — auto-starts on login
Exec=f9-talk --backend cloud
Icon=f9-talk
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=5
Hidden=false
NoDisplay=false
